# hydradragon/Owlyshield/rules/behavior_rules.yaml
- name: "HEUR:Win64.PSW.Stealer.BehavesLike.Exela.gen"
  description: "Detects info stealer behavior targeting browser data and credentials"
  debug: true

  # =============================================================================
  # SINGLE DETECTION SYSTEM — NAMED CONDITIONS
  # =============================================================================

  named_conditions:
    browser_profile_access:
      file_paths:
        - "%USERPROFILE%\\Google\\Chrome\\User Data"
        - "%USERPROFILE%\\AppData\\Local\\Google\\Chrome\\User Data"
        - "Opera Software\\Opera Stable"
        - "Opera Software\\Opera GX Stable"
        - "BraveSoftware\\Brave-Browser\\User Data"
        - "Microsoft\\Edge\\User Data"
      file_operations: ["read", "browse"]
      min_matches: 2

    discord_token_access:
      file_paths:
        - "%APPDATA%\\Discord\\Local Storage"
        - "%APPDATA%\\Discord\\leveldb"
      file_operations: ["read", "browse"]

    credential_file_access:
      file_paths:
        - "Cookies"
        - "Login Data"
        - "Web Data"
        - "Local State"
      file_operations: ["read"]
      min_matches: 2

    crypto_wallet_access:
      file_paths:
        - "*\\Electrum\\wallets*"
        - "*\\Bitcoin\\wallets*"
        - "*\\Ethereum\\keystore*"
      file_operations: ["read"]

    temp_staging:
      file_paths:
        - "%TEMP%"
        - "%APPDATA%\\Local\\Temp"
      file_operations: ["write", "create"]

    decryption_apis:
      apis:
        - "CryptUnprotectData"
        - "BCryptDecrypt"
        - "NCryptDecrypt"
        - "DPAPI"
      min_matches: 1

    archive_creation:
      apis:
        - "CreateZip"
        - "CreateRar"
      process_names:
        - "7z.exe"
        - "tar.exe"
        - "winrar.exe"
      file_extensions:
        - ".zip"
        - ".rar"

    network_exfiltration:
      has_network_activity: true
      network_indicators:
        - "http"
        - "https"
        - "ftp"

    suspicious_origin:
      parent_names:
        - "powershell.exe"
        - "cmd.exe"
        - "wscript.exe"
        - "cscript.exe"
        - "mshta.exe"
        - "python.exe"
        - "node.exe"

    obfuscated_binary:
      entropy_threshold: 7.0

    browser_termination:
      terminated_processes:
        - "msedge.exe"
        - "chrome.exe"
        - "brave.exe"
        - "opera.exe"
        - "firefox.exe"

    startup_persistence:
      registry_keys:
        - "HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run\\*"
      registry_operations: ["set", "create"]

    exfiltration_pattern:
      file_paths:
        - "*Credentials*"
        - "*Passwords*"
        - "*Tokens*"
      file_operations: ["read"]

  # =============================================================================
  # FIXED DETECTION LOGIC - Proper enum structure
  # =============================================================================

  detection_logic:
    # Top-level AND
    and:
      # First AND condition: any_of for data source access
      - any_of:
          - "browser_profile_access"
          - "discord_token_access"
          - "crypto_wallet_access"

      # Second AND condition: any_of for credential access
      - any_of:
          - "credential_file_access"
          - "exfiltration_pattern"

      # Third AND condition: OR for exfiltration method
      - or:
          # Nested AND for archive staging
          - and:
              - condition: "archive_creation"
              - condition: "temp_staging"
          # OR just network exfiltration
          - condition: "network_exfiltration"

      # Fourth AND condition: at_least 2 of these indicators
      - at_least: 2
        conditions:
          - "decryption_apis"
          - "browser_termination"
          - "startup_persistence"
          - "suspicious_origin"
          - "obfuscated_binary"

  # =============================================================================
  # RESPONSE, ALLOWLIST, & METADATA
  # =============================================================================

  response:
    terminate_process: true
    quarantine: true
    auto_revert: true

  allowlisted_apps:
    - pattern: "chrome.exe"
      must_be_signed: true
      signers: ["Google.*"]
    - pattern: "msedge.exe"
      must_be_signed: true
      signers: ["Microsoft.*"]
    - pattern: "brave.exe"
      must_be_signed: true
      signers: ["Brave Software.*"]
    - pattern: "opera.exe"
      must_be_signed: true
      signers: ["Opera.*"]
    - pattern: "discord.exe"
      must_be_signed: true
      signers: ["Discord.*"]
    - pattern: "code.exe"
      must_be_signed: true
      signers: ["Microsoft.*"]
    - pattern: "devenv.exe"
      must_be_signed: true
      signers: ["Microsoft.*"]

  severity: 95
  level: critical
  status: stable
  author: "Emirhan Ucan"

  mitre_attack:
    - "T1555"
    - "T1539"
    - "T1552.001"
    - "T1005"
    - "T1140"
    - "T1560"
    - "T1041"
    - "T1547.001"

  tags:
    - "credential_theft"
    - "browser_stealer"
    - "info_stealer"
    - "data_exfiltration"

  false_positives:
    - "Legitimate password managers"
    - "Browser sync utilities"
    - "Backup software"
    - "Development/debugging tools"
