# hydradragon/Owlyshield/rules/legacy_rules.yaml
- name: "Unified Advanced Stealer Detection"
  description: "Detects multi-stage browser theft, Discord injection, and Exela-style persistence with certificate validation."
  
  # --- Metadata ---
  severity: 95
  level: critical
  status: stable
  mitre_attack:
    - "T1555.003" # Credentials from Web Browsers
    - "T1074.001" # Data Staged
    - "T1560"     # Archive Collected Data
    - "T1574.002" # Hijack Execution Flow (Discord Injection)
    - "T1053.005" # Scheduled Task (Persistence)

  # --- Legacy Stealer Logic (Maintained for Rust line 456 check) ---
  browser_paths:
    - "Google\\Chrome\\User Data"
    - "Opera Software\\Opera Stable"
    - "BraveSoftware\\Brave-Browser\\User Data"
    - "Microsoft\\Edge\\User Data"
    - "Discord\\Local Storage"
  sensitive_files:
    - "Cookies"
    - "Login Data"
    - "Web Data"
    - "Local State"
  staging_paths:
    - "AppData\\Local\\Temp"

  # --- NEW: Stage-Based Behavioral Logic ---
  min_stages_satisfied: 2
  time_window_ms: 60000 

  stages:
    # STAGE 1: Browser & Crypto Harvesting
    - name: "Sensitive_Data_Access"
      conditions:
        - type: "File"
          op: "Read"
          path_pattern: "*\\User Data\\*\\Login Data"
        
        # High Entropy check for encrypted DB reads
        - type: "Heuristic"
          metric: "Entropy"
          threshold: 6.8
          
        # Exela specific: Detects access to multiple Wallet Extension IDs
        - type: "SensitivePathAccess"
          op_type: "Read"
          min_unique_paths: 2
          patterns:
            - "*\\nkbihfbeogaeaoehlefnkodbefgpgknn\\*" # MetaMask
            - "*\\bfnaelmomeimhlpmgjnjophhpkkoljpa\\*" # Phantom
            - "*\\Local Storage\\leveldb\\*.ldb"        # Discord/Electron Tokens

    # STAGE 2: Discord Injection & Modification
    - name: "Client_Integrity_Violation"
      conditions:
        # Exela's Injection method: Writing to the Discord core index.js
        - type: "File"
          op: "Write"
          path_pattern: "*\\modules\\discord_desktop_core\\index.js"
        
        # Killing Discord to unlock files for injection
        - type: "Process"
          op: "CommandLine"
          pattern: "*taskkill* /IM *discord.exe*"

    # STAGE 3: Persistence & Recon (Exela Initialization)
    - name: "Persistence_and_Recon"
      conditions:
        # Detects Exela's specific 'schtasks' persistence and VM checks
        - type: "CommandLineMatch"
          match_mode: "any"
          patterns:
            - pattern: "*schtasks*/create*/tn*UpdateService*" # Exela default task
            - pattern: "*wmic*csproduct*get*uuid*"            # Anti-VM/ID Fingerprint
            - pattern: "*reg*add*CurrentVersion\\Run*"
        
        # File-based persistence
        - type: "File"
          op: "Create"
          path_pattern: "*\\AppData\\Local\\*UpdateService\\*.exe"

    # STAGE 4: Exfiltration & Staging
    - name: "Exfiltration_Attempt"
      conditions:
        - type: "ArchiveCreation"
          extensions: ["zip", "rar", "7z"]
          in_temp: true
          min_size: 20480 # 20KB min to avoid noise
          
        - type: "ByteThreshold"
          direction: "write"
          comparison: "gt"
          threshold: 204800 # 200KB upload/write threshold

  # --- Advanced Certificate-Based Allowlisting ---
  # Only allows these apps if they are signed by the official vendor
  allowlisted_apps:
    - pattern: "explorer.exe"
      must_be_signed: true
      signers: ["Microsoft Windows", "Microsoft Corporation"]
    
    - pattern: "chrome.exe"
      must_be_signed: true
      signers: ["Google LLC"]
    
    - pattern: "msedge.exe"
      must_be_signed: true
      signers: ["Microsoft Corporation"]
    
    - pattern: "firefox.exe"
      must_be_signed: true
      signers: ["Mozilla Corporation", "Mozilla Foundation"]
    
    - pattern: "brave.exe"
      must_be_signed: true
      signers: ["Brave Software, Inc."]
    
    - pattern: "opera.exe"
      must_be_signed: true
      signers: ["Opera Software.*", "Opera Norway AS"]

  # --- Response Actions ---
  response:
    record: true
    terminate_process: true
    quarantine: true
    auto_revert: false
