# hydradragon/Owlyshield/rules/legacy_rules.yaml
- name: "Exela_Stealer_Complete_Detection"
  description: "Comprehensive detection for Exela stealer covering all attack stages: browser data theft, crypto wallets, Discord injection, game session hijacking, telegram sessions, memory scanning, and multi-platform exfiltration"
  
  # --- Metadata ---
  author: "Emirhan Ucan"
  severity: 99
  level: critical
  status: stable
  mitre_attack: 
    - "T1555.003"
    - "T1539"
    - "T1552.001"
    - "T1005"
    - "T1560"
    - "T1041"
    - "T1573"
    - "T1574.002"
    - "T1055"
    - "T1547.001"
    - "T1053.005"
    - "T1140"
    - "T1012"
    - "T1082"
    - "T1518.001"
  
  tags:
    - "stealer"
    - "credential_theft"
    - "browser_hijacking"
    - "discord_injection"
    - "crypto_wallet_theft"
    - "session_hijacking"
    - "telegram_theft"
    - "game_account_theft"
    - "anti_vm"
    - "anti_debug"
    - "persistence"
    - "exfiltration"
  
  # Logic tuning
  require_internet: true
  multi_access_threshold: 3
  entropy_threshold: 6.5
  conditions_percentage: 0.45
  min_stages_satisfied: 3
  time_window_ms: 120000
  
  browser_paths:
    - "*\\Google\\Chrome\\User Data*"
    - "*\\Microsoft\\Edge\\User Data*"
    - "*\\BraveSoftware\\Brave-Browser\\User Data*"
    - "*\\Opera Software\\Opera*"
    - "*\\Mozilla\\Firefox\\Profiles*"
    - "*\\Discord*"
  
  sensitive_files:
    - "Login Data"
    - "Cookies"
    - "Web Data"
    - "Local State"
    - "cookies.sqlite"
    - "places.sqlite"
    - "formhistory.sqlite"
  
  staging_paths:
    - "*\\AppData\\Local\\Temp*"
    - "*\\Temp*"
  
  crypto_apis:
    - "CryptUnprotectData"
  
  suspicious_parents:
    - "cmd.exe"
    - "powershell.exe"
    - "wscript.exe"
    - "mshta.exe"
  
  archive_actions:
    - ".zip"
    - ".rar"
  
  stages:
    # ============================================================================
    # STAGE 1: Anti-Analysis & Evasion
    # ============================================================================
    - name: "Anti_Analysis_Checks"
      conditions:
        - type: "Process"
          op: "CommandLine"
          pattern: "*wmic*csproduct*get*uuid*"
        
        - type: "Process"
          op: "CommandLine"
          pattern: "*wmic*path*win32_VideoController*"
        
        - type: "Process"
          op: "CommandLine"
          pattern: "*taskkill*/F*/IM*wireshark*"
        
        - type: "Process"
          op: "CommandLine"
          pattern: "*taskkill*/F*/IM*processhacker*"
        
        - type: "Process"
          op: "CommandLine"
          pattern: "*tasklist*"
        
        - type: "File"
          op: "Read"
          path_pattern: "D:\\Tools*"

    # ============================================================================
    # STAGE 2: Browser Data Harvesting
    # ============================================================================
    - name: "Browser_Credential_Theft"
      conditions:
        - type: "SensitivePathAccess"
          op_type: "Read"
          min_unique_paths: 3
          patterns:
            - "*\\User Data\\*\\Login Data"
            - "*\\User Data\\*\\Cookies"
            - "*\\User Data\\*\\Web Data"
            - "*\\User Data\\Local State"
            - "*\\Network\\Cookies"
        
        - type: "File"
          op: "Read"
          path_pattern: "*\\Firefox\\Profiles\\*\\cookies.sqlite"
        
        - type: "File"
          op: "Read"
          path_pattern: "*\\Firefox\\Profiles\\*\\places.sqlite"
        
        - type: "Process"
          op: "CommandLine"
          pattern: "*taskkill*/F*/IM*chrome.exe*"
        
        - type: "Process"
          op: "CommandLine"
          pattern: "*taskkill*/F*/IM*firefox.exe*"
        
        - type: "File"
          op: "Create"
          path_pattern: "*\\Temp\\*.db"
        
        - type: "Heuristic"
          metric: "Entropy"
          threshold: 6.5

    # ============================================================================
    # STAGE 3: Cryptocurrency Wallet Theft
    # ============================================================================
    - name: "Crypto_Wallet_Harvesting"
      conditions:
        - type: "SensitivePathAccess"
          op_type: "Read"
          min_unique_paths: 2
          patterns:
            - "*\\Local Extension Settings\\nkbihfbeogaeaoehlefnkodbefgpgknn\\*"
            - "*\\Local Extension Settings\\fhbohimaelbohpjbbldcngcnapndodjp\\*"
            - "*\\Local Extension Settings\\bfnaelmomeimhlpmgjnjophhpkkoljpa\\*"
            - "*\\Local Extension Settings\\hnfanknocfeofbddgcijnmhnfnkdnaad\\*"
            - "*\\Local Extension Settings\\fnjhmkhhmkbjkkabndcnnogagogbneec\\*"
            - "*\\Local Extension Settings\\aholpfdialjgjfhomihkjbmgjidlcdno\\*"
            - "*\\Bitcoin\\wallets\\*"
            - "*\\Exodus\\exodus.wallet\\*"
            - "*\\Ethereum\\keystore\\*"
            - "*\\Electrum\\wallets\\*"
        
        - type: "File"
          op: "Write"
          path_pattern: "*\\Temp\\*\\Wallets\\*"

    # ============================================================================
    # STAGE 4: Discord Token & Injection
    # ============================================================================
    - name: "Discord_Token_Theft_And_Injection"
      conditions:
        - type: "SensitivePathAccess"
          op_type: "Read"
          min_unique_paths: 1
          patterns:
            - "*\\discord\\Local Storage\\leveldb\\*.ldb"
            - "*\\discord\\Local Storage\\leveldb\\*.log"
            - "*\\discordcanary\\Local Storage\\leveldb\\*.ldb"
            - "*\\discordptb\\Local Storage\\leveldb\\*.ldb"
        
        - type: "MemoryScan"
          patterns:
            - "dQw4w9WgXcQ:"
          detect_pe_headers: false
          private_only: false
        
        - type: "File"
          op: "Write"
          path_pattern: "*\\discord\\app-*\\modules\\discord_desktop_core*\\index.js"
        
        - type: "Process"
          op: "CommandLine"
          pattern: "*taskkill*/F*/IM*Discord.exe*"
        
        - type: "Process"
          op: "CommandLine"
          pattern: "*Update.exe*--processStart*Discord.exe*"
        
        - type: "Network"
          op: "Connect"
          dest_pattern: "*discord.com/api/v*/users/@me*"

    # ============================================================================
    # STAGE 5: Social Media Session Hijacking
    # ============================================================================
    - name: "Social_Media_Session_Theft"
      conditions:
        - type: "Network"
          op: "Connect"
          dest_pattern: "*i.instagram.com/api/v1/accounts/current_user*"
        
        - type: "Network"
          op: "Connect"
          dest_pattern: "*tiktok.com/passport/web/account/info*"
        
        - type: "Network"
          op: "Connect"
          dest_pattern: "*twitter.com/i/api/*/account/update_profile*"
        
        - type: "Network"
          op: "Connect"
          dest_pattern: "*oauth.reddit.com/api/v1/me*"
        
        - type: "Network"
          op: "Connect"
          dest_pattern: "*spotify.com/api/account-settings*"
        
        - type: "Network"
          op: "Connect"
          dest_pattern: "*roblox.com/my/account/json*"
        
        - type: "Network"
          op: "Connect"
          dest_pattern: "*gql.twitch.tv/gql*"

    # ============================================================================
    # STAGE 6: Gaming Platform Theft
    # ============================================================================
    - name: "Gaming_Session_Hijacking"
      conditions:
        - type: "File"
          op: "Read"
          path_pattern: "C:\\Program Files (x86)\\Steam\\config\\loginusers.vdf"
        
        - type: "Network"
          op: "Connect"
          dest_pattern: "*api.steampowered.com/ISteamUser/GetPlayerSummaries*"
        
        - type: "File"
          op: "Write"
          path_pattern: "*\\Temp\\*\\Games\\Steam\\*"
        
        - type: "File"
          op: "Read"
          path_pattern: "*\\Ubisoft Game Launcher\\*"
        
        - type: "File"
          op: "Read"
          path_pattern: "*\\EpicGamesLauncher\\Saved\\Config\\Windows\\*"
        
        - type: "File"
          op: "Read"
          path_pattern: "*\\Growtopia\\save.dat"

    # ============================================================================
    # STAGE 7: Telegram Session Theft
    # ============================================================================
    - name: "Telegram_Session_Hijacking"
      conditions:
        - type: "Process"
          op: "CommandLine"
          pattern: "*taskkill*/F*/IM*Telegram.exe*"
        
        - type: "SensitivePathAccess"
          op_type: "Read"
          min_unique_paths: 2
          patterns:
            - "*\\Telegram Desktop\\tdata\\*"
            - "*\\Telegram Desktop\\tdata\\D877F783D5D3EF8C*"
        
        - type: "File"
          op: "Write"
          path_pattern: "*\\Temp\\*\\Telegram Session\\*"

    # ============================================================================
    # STAGE 8: System Reconnaissance
    # ============================================================================
    - name: "System_Reconnaissance"
      conditions:
        - type: "Process"
          op: "CommandLine"
          pattern: "*systeminfo*"
        
        - type: "Process"
          op: "CommandLine"
          pattern: "*hostname*"
        
        - type: "Process"
          op: "CommandLine"
          pattern: "*wmic*logicaldisk*get*"
        
        - type: "Process"
          op: "CommandLine"
          pattern: "*net*user*"
        
        - type: "Process"
          op: "CommandLine"
          pattern: "*ipconfig*/all*"
        
        - type: "Process"
          op: "CommandLine"
          pattern: "*netstat*-ano*"
        
        - type: "Process"
          op: "CommandLine"
          pattern: "*netsh*wlan*show*profiles*"
        
        - type: "Process"
          op: "CommandLine"
          pattern: "*netsh*wlan*show*profile*key=clear*"
        
        - type: "Network"
          op: "Connect"
          dest_pattern: "*ip-api.com/json*"
        
        - type: "Process"
          op: "CommandLine"
          pattern: "*tasklist*/FO*LIST*"
        
        - type: "Process"
          op: "CommandLine"
          pattern: "*powershell*Get-Clipboard*"

    # ============================================================================
    # STAGE 9: Screenshot Capture
    # ============================================================================
    - name: "Screen_Capture"
      conditions:
        - type: "Process"
          op: "CommandLine"
          pattern: "*powershell*-EncodedCommand*"
        
        - type: "Process"
          op: "CommandLine"
          pattern: "*System.Drawing*"
        
        - type: "Process"
          op: "CommandLine"
          pattern: "*CopyFromScreen*"
        
        - type: "File"
          op: "Create"
          path_pattern: "*\\Display (*).png"

    # ============================================================================
    # STAGE 10: Data Staging & Archiving
    # ============================================================================
    - name: "Data_Exfiltration_Staging"
      conditions:
        - type: "File"
          op: "Create"
          path_pattern: "*\\Temp\\*\\Browsers\\*"
        
        - type: "File"
          op: "Create"
          path_pattern: "*\\Temp\\*\\Sessions\\*"
        
        - type: "File"
          op: "Create"
          path_pattern: "*\\Temp\\*\\Tokens\\*"
        
        - type: "File"
          op: "Create"
          path_pattern: "*\\Temp\\*\\Wallets\\*"
        
        - type: "File"
          op: "Write"
          path_pattern: "*\\Temp\\*\\Passwords.txt"
        
        - type: "File"
          op: "Write"
          path_pattern: "*\\Temp\\*\\Cookies.txt"
        
        - type: "ArchiveCreation"
          extensions: ["zip"]
          in_temp: true
          min_size: 50000
        
        - type: "ByteThreshold"
          direction: "write"
          comparison: "gt"
          threshold: 102400

    # ============================================================================
    # STAGE 11: Exfiltration via Discord Webhook
    # ============================================================================
    - name: "Webhook_Exfiltration"
      conditions:
        - type: "Network"
          op: "Connect"
          dest_pattern: "*discord.com/api/webhooks/*"
        
        - type: "Network"
          op: "Connect"
          dest_pattern: "*discordapp.com/api/webhooks/*"
        
        - type: "MemoryScan"
          patterns:
            - "https://discord.com/api/webhooks/"
            - "https://discordapp.com/api/webhooks/"
            - "webhook"
          detect_pe_headers: false
          private_only: false
        
        - type: "Network"
          op: "Connect"
          dest_pattern: "*api.gofile.io/getServer*"
        
        - type: "Network"
          op: "Connect"
          dest_pattern: "*gofile.io/uploadFile*"
        
        - type: "MemoryScan"
          patterns:
            - "Exela Stealer"
            - "t.me/ExelaStealer"
          detect_pe_headers: false
          private_only: false

    # ============================================================================
    # STAGE 12: Persistence Mechanisms
    # ============================================================================
    - name: "Persistence_Installation"
      conditions:
        - type: "File"
          op: "Create"
          path_pattern: "*\\AppData\\Local\\*UpdateService\\*.exe"
        
        - type: "File"
          op: "Create"
          path_pattern: "*\\AppData\\Local\\ExelaUpdateService\\Exela.exe"
        
        - type: "Process"
          op: "CommandLine"
          pattern: "*attrib*+h*+s*UpdateService*"
        
        - type: "Process"
          op: "CommandLine"
          pattern: "*schtasks*/create*/tn*ExelaUpdateService*"
        
        - type: "Process"
          op: "CommandLine"
          pattern: "*schtasks*/create*/sc*onlogon*/rl*highest*"
        
        - type: "Registry"
          op: "SetValue"
          key_pattern: "*\\Software\\Microsoft\\Windows\\CurrentVersion\\Run*"
          value_name: "Exela Update Service"
        
        - type: "Process"
          op: "CommandLine"
          pattern: "*reg*add*CurrentVersion\\Run*/v*Exela Update Service*"
        
        - type: "File"
          op: "Create"
          path_pattern: "*\\Start Menu\\Programs\\Startup\\Exela.exe"

    # ============================================================================
    # STAGE 13: Document And File Theft
    # ============================================================================
    - name: "Document_And_File_Theft"
      conditions:
        - type: "SensitivePathAccess"
          op_type: "Read"
          min_unique_paths: 3
          patterns:
            - "*\\Desktop\\*.txt"
            - "*\\Desktop\\*.pdf"
            - "*\\Desktop\\*.doc*"
            - "*\\Documents\\*.txt"
            - "*\\Documents\\*.pdf"
            - "*\\Downloads\\*.txt"
        
        - type: "File"
          op: "Read"
          path_pattern: "*secret*"
        
        - type: "File"
          op: "Read"
          path_pattern: "*password*"
        
        - type: "File"
          op: "Read"
          path_pattern: "*wallet*"
        
        - type: "File"
          op: "Write"
          path_pattern: "*\\Temp\\StealedFilesByExela\\*"

    # ============================================================================
    # STAGE 14: Social Engineering Distraction
    # ============================================================================
    - name: "Social_Engineering_Distraction"
      conditions:
        - type: "Process"
          op: "CommandLine"
          pattern: "*mshta*javascript*WScript.Shell*Popup*"
        
        - type: "Process"
          op: "CommandLine"
          pattern: "*api-ms-win-crt-runtime*missing*"
        
        - type: "Process"
          op: "Spawn"
          pattern: "*mshta.exe*"

  # --- Correlation Logic ---
  mapping:
    or:
      - and:
          - stage: "Browser_Credential_Theft"
          - stage: "Data_Exfiltration_Staging"
          - stage: "Webhook_Exfiltration"
      - and:
          - stage: "Crypto_Wallet_Harvesting"
          - stage: "Webhook_Exfiltration"
      - and:
          - stage: "Discord_Token_Theft_And_Injection"
          - stage: "Webhook_Exfiltration"
      - and:
          - stage: "Gaming_Session_Hijacking"
          - stage: "Data_Exfiltration_Staging"
      - and:
          - stage: "Telegram_Session_Hijacking"
          - stage: "Data_Exfiltration_Staging"
      - and:
          - stage: "Anti_Analysis_Checks"
          - stage: "Persistence_Installation"
          - stage: "Webhook_Exfiltration"

  # --- Memory Scan Configuration ---
  memory_scan_config:
    target_processes: ["python*", "pythonw*", "*.exe"]
    scan_on_io_event: false
    scan_every_n_ops: 250

  # --- Certificate Allowlisting ---
  allowlisted_apps:
    - pattern: "explorer.exe"
      must_be_signed: true
      signers: ["Microsoft Windows", "Microsoft Corporation"]
    
    - pattern: "chrome.exe"
      must_be_signed: true
      signers: ["Google LLC"]
    
    - pattern: "msedge.exe"
      must_be_signed: true
      signers: ["Microsoft Corporation"]
    
    - pattern: "firefox.exe"
      must_be_signed: true
      signers: ["Mozilla Corporation", "Mozilla Foundation"]
    
    - pattern: "brave.exe"
      must_be_signed: true
      signers: ["Brave Software, Inc."]

  # --- Automated Response ---
  response:
    record: true
    terminate_process: true
    quarantine: true
    suspend_process: false
    auto_revert: false
  
  proximity_log_threshold: 70.0
  debug: false
